# etc/supervisor/conf.d/exiter.conf
# Demo setup for script that exits with zero status (exit status by default)

[program:exiter]
command=/usr/bin/php -f /home/user/exiter.php
user=user
