# etc/supervisor/conf.d/failer.conf
# Demo setup for script that fails with non-zero status and will be restarted

[program:failer]
command=/usr/bin/php -f /home/user/failer.php
user=user
