# etc/supervisor/conf.d/myother.conf
# Process setup that will start two processes under user user

[program:myother]
command=/usr/bin/php -f /home/user/other.php
user=user
numprocs=2
process_name=%(program_name)s_%(process_num)s
