# etc/supervisor/conf.d/myprocess.conf
# Simplest process setup, process will be started as root

[program:myprocess]
command=/usr/bin/php -f /home/user/long.php
